# 10.4.3/build/rootfs/etc/profile.d/env.sh
#!/bin/bash

MYSQL_HOST="172.42.0.1"
export MYSQL_HOST

